# contrib/init.d/jack-bluealsa-autobridge
#!/bin/sh
### BEGIN INIT INFO
# Provides:          jack-bluealsa-autobridge
# Required-Start:    $local_fs $remote_fs dbus bluetooth bluealsad jackd-rt
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Autobridge BlueALSA PCMs into JACK via alsa_in/alsa_out
# Description:       Starts jack-bluealsa-autobridge after JACK and BlueALSA are ready.
### END INIT INFO

. /lib/lsb/init-functions

# Load overrides for user detection from jackd-rt defaults (optional)
if [ -r /etc/default/jackd-rt ]; then
    . /etc/default/jackd-rt
fi

DAEMON="/usr/local/bin/jack-bluealsa-autobridge"
PIDFILE="/var/run/jack-bluealsa-autobridge.pid"
LOGFILE="/var/log/jack-bluealsa-autobridge.log"

start() {
    log_daemon_msg "Starting jack-bluealsa-autobridge" "jack-bluealsa-autobridge"

    if [ ! -x "$DAEMON" ]; then
        log_failure_msg "Binary not found: $DAEMON"
        log_end_msg 1
        return 1
    fi

    # Detect user to run under (match jackd-rt logic)
    USERNAME=""
    if [ -n "$JACKD_USER" ]; then
        USERNAME="$JACKD_USER"
    else
        # first user with UID>=1000 and in audio group
        for u in $(cut -d: -f1 /etc/passwd); do
            uid=$(id -u "$u" 2>/dev/null || true)
            if [ -n "$uid" ] && [ "$uid" -ge 1000 ]; then
                groups=$(id -nG "$u" 2>/dev/null || true)
                if printf '%s\n' " $groups " | grep -q ' audio '; then
                    USERNAME="$u"; break
                fi
            fi
        done
    fi
    if [ -z "$USERNAME" ]; then
        # fallback to first UID>=1000
        for u in $(cut -d: -f1 /etc/passwd); do
            uid=$(id -u "$u" 2>/dev/null || true)
            if [ -n "$uid" ] && [ "$uid" -ge 1000 ]; then
                USERNAME="$u"; break
            fi
        done
    fi
    [ -z "$USERNAME" ] && USERNAME="root"

    # Wait for JACK to be ready
    if command -v jack_wait >/dev/null 2>&1; then
        jack_wait -w -t 15 >/dev/null 2>&1 || true
    else
        sleep 2
    fi

    # Wait for BlueALSA D-Bus service org.bluealsa (up to 15s)
    i=0
    while [ $i -lt 15 ]; do
        if command -v gdbus >/dev/null 2>&1; then
            if gdbus call --system --dest org.freedesktop.DBus \
                --object-path /org/freedesktop/DBus \
                --method org.freedesktop.DBus.NameHasOwner org.bluealsa 2>/dev/null | grep -q true; then
                break
            fi
        fi
        sleep 1
        i=$((i+1))
    done

    touch "$LOGFILE" 2>/dev/null || true
    chown "$USERNAME":"$USERNAME" "$LOGFILE" 2>/dev/null || true

    CMD="exec su -l \"$USERNAME\" -s /bin/sh -c 'exec \"$DAEMON\" >>\"$LOGFILE\" 2>&1'"

    start-stop-daemon --start --background --make-pidfile --pidfile "$PIDFILE" \
        --startas /bin/sh -- -c "$CMD" \
        || { log_failure_msg "Failed to start jack-bluealsa-autobridge"; log_end_msg 1; return 1; }

    log_end_msg 0
    return 0
}

stop() {
    log_daemon_msg "Stopping jack-bluealsa-autobridge" "jack-bluealsa-autobridge"
    start-stop-daemon --stop --pidfile "$PIDFILE" --retry 5 || true
    rm -f "$PIDFILE" 2>/dev/null || true
    log_end_msg 0
}

status() {
    if [ -f "$PIDFILE" ]; then
        pid=$(cat "$PIDFILE" 2>/dev/null)
        if [ -n "$pid" ] && kill -0 "$pid" >/dev/null 2>&1; then
            echo "jack-bluealsa-autobridge running (pid $pid)"
            return 0
        fi
    fi
    echo "jack-bluealsa-autobridge not running"
    return 1
}

case "$1" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    restart|force-reload)
        stop
        sleep 1
        start
        ;;
    status)
        status
        ;;
    *)
        echo "Usage: $0 {start|stop|restart|status}"
        exit 2
        ;;
esac

exit 0