# contrib/init.d/jack-bridge-loopback
#!/bin/sh
### BEGIN INIT INFO
# Provides:          jack-bridge-loopback
# Required-Start:    $local_fs $remote_fs
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Ensure ALSA snd-aloop Loopback card is available for Gaming
# Description:       Loads snd-aloop at boot (if available) so pcm.gaming and
#                    the gaming_in JACK bridge work reliably.
### END INIT INFO

. /lib/lsb/init-functions

ensure_loopback_module() {
    # If Loopback card already exists, nothing to do.
    if [ -r /proc/asound/cards ] && \
       grep -q 'Loopback' /proc/asound/cards 2>/dev/null; then
        return 0
    fi

    # If modprobe is available, try to load snd-aloop.
    if command -v modprobe >/dev/null 2>&1; then
        if modprobe snd-aloop 2>/dev/null; then
            # Give kernel a moment to register the card
            sleep 1
            return 0
        else
            log_warning_msg "jack-bridge-loopback: modprobe snd-aloop failed; Gaming bridge will not be available"
            return 1
        fi
    fi

    # No modprobe and no Loopback: log once and exit.
    log_warning_msg "jack-bridge-loopback: snd-aloop not available (no Loopback card, no modprobe)"
    return 1
}

case "$1" in
  start)
    log_daemon_msg "Ensuring ALSA Loopback card (snd-aloop)" "jack-bridge-loopback"
    ensure_loopback_module
    log_end_msg 0
    ;;
  stop)
    # No-op: we do not unload snd-aloop; leaving Loopback available is harmless.
    ;;
  restart|force-reload)
    # Nothing to restart; module presence is static for the life of the boot.
    ;;
  status)
    if [ -r /proc/asound/cards ] && grep -q 'Loopback' /proc/asound/cards 2>/dev/null; then
        echo "Loopback card present"
        exit 0
    else
        echo "Loopback card not present"
        exit 1
    fi
    ;;
  *)
    echo "Usage: $0 {start|stop|restart|status}"
    exit 2
    ;;
esac

exit 0

